# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/aaastatsd.service
[Unit]
Description=AAA Statistics Collection daemon
Requires=hostcfgd.service
After=updategraph.service interfaces-config.service hostcfgd.service

[Service]
Type=simple
ExecStart=/usr/bin/aaastatsd

[Install]
WantedBy=multi-user.target
