# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/critical-monitoring.service
[Unit]
Description=Critical Memory Threshold Monitoring Service
After=rc-local.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/critical-monitoring.sh
KillSignal=SIGKILL
SuccessExitStatus=SIGKILL

[Install]
WantedBy=multi-user.target
