# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/database.service
[Unit]
Description=Database container

Wants=database-chassis.service
After=database-chassis.service
Requires=docker.service
After=docker.service
After=rc-local.service
StartLimitIntervalSec=1200
StartLimitBurst=3

[Service]
User=root
ExecStartPre=/usr/bin/database.sh start
ExecStart=/usr/bin/database.sh wait
ExecStop=/usr/bin/database.sh stop
Restart=always
RestartSec=30
Slice=sonic-core.slice

[Install]
WantedBy=multi-user.target
