# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/histogram.service
[Unit]
Description=Memory Histogram
Before=rc-local.service

[Service]
Type=forking
StandardOutput=syslog+console
ExecStart=/usr/local/bin/histogramd.py --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
