# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/in-memory-log-rotate-daemon.service
[Unit]
Description=In-memory auto log rotation infrastructure
Before=rc-local.service

[Service]
Type=simple
StandardOutput=syslog
ExecStart=/usr/bin/in-memory-log-rotate-daemon.sh start 
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
