# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/sysmonitor.service
[Unit]
Description=Watchdog service for system services
Requires=database.service
After=database.service
StartLimitIntervalSec=60
StartLimitBurst=3


[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/sysmonitor.py
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
