# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/warmboot-finalizer.service
[Unit]
Description=Monitor warm recovery and disable warmboot when done
Requires=database.service
After=database.service

[Service]
Type=simple
ExecStart=/usr/local/bin/finalize-warmboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
