# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/watchdog-control.service
[Unit]
Description=watchdog control service
After=rc-local.service platform-ready.service

[Service]
Type=forking
Restart=always
ExecStart=/usr/local/bin/watchdog-control.sh start 
ExecStop=/usr/local/bin/watchdog-control.sh stop
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
