# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/multi-user.target.wants/ztp-config.service
[Unit]
Description=ZTP configuration service
Requires=database.service
After=database.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ztp-config
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
