# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/tpcm@mgmt.service.d/dependency.conf
[Unit]
After=network-online.target hostcfgd.service
Requires=network-online.target hostcfgd.service
