# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/sonic.target.wants/autobreakoutmgrd.service
[Unit]
Description=Auto Breakout manager daemon
Requires=database.service updategraph.service 
After=updategraph.service database.service pmon.service
BindsTo=sonic.target
After=sonic.target
StartLimitIntervalSec=1200
StartLimitBurst=20

[Service]
Type=simple
ExecStart=/usr/local/bin/autobreakoutmgrd
Restart=always
RestartSec=30

[Install]
WantedBy=sonic.target
