# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/sonic.target.wants/lacp_helper.service
[Unit]
Description=LACP Helper service
Requires=database.service
After=database.service
BindsTo=sonic.target
After=sonic.target
StartLimitIntervalSec=1200
StartLimitBurst=3

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lacp_helper.py

[Install]
WantedBy=sonic.target
