# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/sonic.target.wants/sonic-init-updatedb.service
[Unit]
Description=Retrieve software information from files and store them into StateDB
Requires=database.service
After=database.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sonic-init-updatedb
Restart=on-failure
RestartSec=10
RemainAfterExit=yes

[Install]
WantedBy=sonic.target
