# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/rc1.d/K01ndctl-monitor
#!/lib/init/init-d-script
### BEGIN INIT INFO
# Provides:             ndctl-monitor
# Required-Start:       $syslog $time $remote_fs
# Required-Stop:        $syslog $time $remote_fs
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description: Ndctl Monitor Daemon
# Description: monitoring of smart events of nvdimm objects
### END INIT INFO
DAEMON=/usr/bin/ndctl
DAEMON_ARGS="monitor --daemon"
