# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/ssh.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/local/bin/host-ssh-setup.sh
ExecStartPre=/usr/sbin/sshd -t
StartLimitIntervalSec=120
StartLimitBurst=25
