# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/systemd/system/tpcm.service
[Unit]
Description=TPCM container
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tpcm.sh start
ExecStop=/usr/local/bin/tpcm.sh stop
StandardOutput=syslog
StandardError=syslog
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
