# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/cron.d/core_cleanup
# Attempts to clean up core files every 10 min
*/10 * * * * root /usr/bin/core_cleanup.py > /dev/null 2>&1

