# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/cron.d/logrotate
# Attempt to rotate system logs once every 5 minutes.
# This wrapper will logrotate as well as cleans old archive files
*/5 * * * * root /usr/bin/log-rotate.sh > /dev/null 2>&1
