# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/cron.d/mem-histogram-rotate
# Attempt to rotate forcefully the /var/log/mem-histogram.log file for every day
58 23 * * * root /usr/sbin/logrotate -f /etc/logrotate.d/histogram/mem-histogram-rule  > /dev/null 2>&1
