# extracted_showtech/sonic_dump_MYBLE01-DC-0409-SI02_20260316_055247/etc/cron.d/run_fetchcrl
#Run fetch-crl-utility once in every 30 minutes.
*/30 * * * * root /usr/bin/python3 /usr/bin/fetch-crl-utility.py > /dev/null 2>&1
